# Desktop/Agent_System/inbox-daemon.service
[Unit]
Description=Inbox Auto-Processing Daemon — watches Input/ for intent files
After=network.target

[Service]
Type=simple
User=frost
ExecStart=/usr/bin/python3 /home/frost/Desktop/Agent_System/inbox_daemon.py --poll 3
Restart=on-failure
RestartSec=5
Environment=INBOX_POLL_INTERVAL=3.0
Environment=INBOX_DEBOUNCE=2.0

[Install]
WantedBy=multi-user.target
