# Desktop/Agent_System/daemons/consistency-daemon/uncertainty-organ.service
[Unit]
Description=Uncertainty Organ — Live confidence map and pre-generation halt gate
After=network.target sovereign.service

[Service]
Type=simple
User=frost
WorkingDirectory=/home/frost/Desktop/Agent_System/consistency-daemon
ExecStart=/home/frost/Desktop/Agent_System/consistency-daemon/.venv/bin/python3 -m consistency_daemon.uncertainty_organ
Restart=on-failure
RestartSec=5
Environment=INBOX_CLARIFY_DIR=/home/frost/Desktop/Inbox/Clarify
Environment=UNCERTAINTY_THRESHOLD=0.70
Environment=UNCERTAINTY_ORGAN_PORT=8439

[Install]
WantedBy=multi-user.target
